# sparse/usr/lib/systemd/user/jolla-actdead-charging.service.d/50-compositor.conf
[Service]
# All devices: Make previous UI exit
# When applicable: Have freshly started android hwc available
ExecStartPre=/usr/sbin/dummy_compositor --hwc-restart --exit-on-enable
