# sparse-10/usr/lib/systemd/system/sailfish-upgrade-ui.service.d/50-vendor.hwcomposer-2-3.conf
[Service]
# make sure unlock-ui has exited before starting upgrade-ui
ExecStartPre=-/usr/sbin/dummy_compositor --exit-on-enable
