# sparse-10/usr/lib/systemd/user/lipstick.service.d/50-vendor.hwcomposer-2-3.conf
[Service]
# make unlock ui exit
ExecStartPre=/usr/sbin/dummy_compositor
