# sparse-10/usr/lib/systemd/system/sailfish-unlock-agent.service.d/50-vendor.hwcomposer-2-3.conf
[Service]
# stop hwcomposer before unlock ui
ExecStartPre=-/system/bin/stop vendor.hwcomposer-2-3

# start hwcomposer after unlock ui, but not on upgrade target
ExecStart=-/bin/sh -c '/usr/bin/test -f /tmp/os-update-running || /system/bin/start vendor.hwcomposer-2-3'
